# backend/crontab
30 */2 * * * cd /app && /venv/bin/python -m src.daemon > /app/daemon/daemon.log 2>&1
