# backend/docker/crontab
*/5 * * * * /app/output/daemon > /app/daemon/daemon.log 2>&1
